# cross/aarch64_be-unknown-linux-gnu/1-base/root/etc/profile.d/makeopts.sh
#!/bin/bash

export MAKEOPTS="-j$(lscpu -p | grep -cv '#')"
